# Network/pi2_wlan1/udev/rules.d/11-openplotter-usb0.rules
KERNEL=="usb0", SUBSYSTEMS=="net", RUN+="/bin/bash /home/pi/.config/openplotter/Network/11-openplotter-usb0.sh"
